# docker/php/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/sbin/apache2ctl -DFOREGROUND

[program:php-fpm]
command=/usr/sbin/php-fpm5.6 --nodaemonize
